# Docker/visorblog.conf
[program:fancoblog]
command=pipenv run gunicorn -w 4 wsgi:app
directory=/usr/src/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true